# deploy/renty-api.service
[Unit]
Description=Renty Dynamic Pricing API Server
After=network.target

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/opt/renty-pricing
Environment="PATH=/opt/renty-pricing/venv/bin"
ExecStart=/opt/renty-pricing/venv/bin/uvicorn api_server:app --host 0.0.0.0 --port 8000
Restart=always
RestartSec=5

# Logging
StandardOutput=append:/var/log/renty/api.log
StandardError=append:/var/log/renty/api-error.log

[Install]
WantedBy=multi-user.target

